# NoSPAM/Install/root/service/dcc/run
#!/bin/sh

exec 2>&1 \
envdir ./env \
sh -c ' exec \
	/usr/bin/dccifd -b
'
