# NoSPAM/Install/root/etc/cron.hourly/ntpdate
#!/bin/sh

ntpdate time.windows.com time.nist.gov> /dev/null 2>&1
