# NoSPAM/Install/root/service/named/run
#!/bin/sh

exec 2>&1 \
sh -c ' exec /usr/sbin/named -u named -f'
