# NoSPAM/Install/root/service/pyzor/run
#!/bin/sh

exec 2>&1 \
envdir ./env \
sh -c ' exec \
	/usr/bin/readyexecd.py /home/NoSPAM/sa/pyzor/readyexec pyzor.client.run
'
