# NoSPAM/Install/root/service/gad/run
#!/bin/sh

exec 2>&1 \
envdir ./env \
sh -c ' exec \
	/home/NoSPAM/bin/ga-daemon
'
