# NoSPAM/Install/root/service/gad/run
#!/bin/sh

exec /home/NoSPAM/bin/ga-daemon
