# NoSPAM/Install/root/service/nsd/run
#!/bin/sh

killall ns-daemon_ex
sleep 1
killall -9 ns-daemon_ex

exec 2>&1 \
envdir ./env \
sh -c ' exec \
	/home/NoSPAM/bin/ns-daemon_ex
'
