# NoSPAM/Install/root/service/nsd/run
#!/bin/sh

exec /home/NoSPAM/bin/ns-daemon
